# recipes-support/ipmitool/files/89-ipmi-i2c.rules
SUBSYSTEM=="module", ACTION=="add", KERNEL=="ipmi_devintf", RUN+="/etc/udev/scripts/ipmi.init.dev"
